# DEBIAN/prerm
#!/bin/sh

systemctl disable spacy

systemctl stop spacy

exit 0
